# woos-pinning.service
[Unit]
Description=woos A-type pinning scanner
After=network-online.target

[Service]
Type=simple
User=neosiwon
EnvironmentFile=/home/neosiwon/woos-alpha-bot/.env
ExecStart=/usr/bin/python3 /home/neosiwon/woos_pinning_scanner_v3.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
